# gitfetch.timer
[Unit]
Description=Git fetch timer

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
