# gitfetch.service
[Unit]
Description=Git fetch
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/gitfetch
